# screendumpLowFrame/layout/DEBIAN/postinst
#!/bin/sh

if [ -L "/var/jb" ]; then
    launchctl load /var/jb/Library/LaunchDaemons/ru.mostmodest.screendumpd.plist
else
    launchctl load /Library/LaunchDaemons/ru.mostmodest.screendumpd.plist
fi

exit 0;
